# configs/aiohttp.conf
[program:aiohttp]
command=/weather/.env/bin/python server.py
directory=/weather/py_aiohttp/
user=nobody
autorestart=true
redirect_stderr=true