# configs/bottle.conf
[program:bottle]
command=/weather/.env/bin/python server.py
directory=/weather/py_bottle/
user=nobody
autorestart=true
redirect_stderr=true