# configs/flask.conf
[program:flask]
command=/weather/.env/bin/python server.py
directory=/weather/py_flask/
user=nobody
autorestart=true
redirect_stderr=true