# configs/tornado.conf
[program:tornado]
command=/weather/.env/bin/python server.py
directory=/weather/py_tornado/
user=nobody
autorestart=true
redirect_stderr=true