# configs/wsgiref.conf
[program:wsgiref]git add .
command=/weather/.env/bin/python server.py
directory=/weather/py_wsgiref/
user=nobody
autorestart=true
redirect_stderr=true